# puppet-modules/cloudaware/files/cron.d/facts2ca
MAILTO=""
PATH=/opt/puppet/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
*/15 * * * * root mco facts2ca
