# ansible-to-s3/etc/cron.d/ansible_to_s3
MAILTO=""
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
*/15 * * * * root ansible-playbook /etc/ansible/facts_to_s3.yaml
