# ohai-to-s3/etc/cron.d/ohai2s3
MAILTO=""
PATH=/opt/puppet/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
*/15 * * * * /usr/sbin/ohai2s3.rb
